# zme_radiotools.crontab.in
#MAILTO=""

0 1  * * *  root	@ZME_RTHOME@/zme_logsclean.sh 30
